# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:webserver]
command=/bin/bash -c "cd /home && python -m SimpleHTTPServer 80"